# etc/sensor-monitor.service
[Unit]
Description=Monitor Sensors & Web Services
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/tmp
ExecStart=/usr/local/bin/sensor-monitor
Restart=on-abort

[Install]
WantedBy=multi-user.target
